# srcs/ftps/supervisord.conf
[supervisord]
user = root
nodaemon = true

[program:vsftpd]
command = sh -c "/usr/sbin/sshd -D && kill -s SIGTERM $(cat supeervisord.pid)"